# deploy/germany-server/bgutil-pot.service
# /etc/systemd/system/bgutil-pot.service
# Exposes bgutil POT to Render (and other clients) on port 26148.
#
# Apply on Germany server:
#   sudo bash deploy/germany-server/configure-bgutil-pot-external.sh

[Unit]
Description=SafeTube bgutil POT HTTP provider (external, port 26148)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/bgutil-pot server --host 0.0.0.0 --port 26148
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
