# deploy/germany-server/bgutil-pot.service.example
# /etc/systemd/system/bgutil-pot.service
# Rust POT provider (jim60105/bgutil-ytdlp-pot-provider-rs). Install binary as /usr/local/bin/bgutil-pot
# and ensure yt-dlp plugin is under the bridge server/yt-dlp-plugins/ directory.
#
#   sudo cp bgutil-pot.service.example /etc/systemd/system/bgutil-pot.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now bgutil-pot.service
#   curl -sS http://127.0.0.1:26148/ping
#
# For Render (external access), prefer:
#   sudo bash deploy/germany-server/configure-bgutil-pot-external.sh

[Unit]
Description=SafeTube bgutil POT HTTP provider for yt-dlp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/bgutil-pot server --host 0.0.0.0 --port 26148
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
